# dosbox.conf
[program:dosbox]
command=dosbox -c "mount c /root/dos" -c "c:" -c "BWAVE386.EXE"
autorestart=true

